# samba-dc/scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:samba]
command=/usr/sbin/samba -i

